# packages/nmosreverseproxy/nmosreverseproxy-0.6.0.tar.gz/nmosreverseproxy-0.6.0/debian/python3-nmosreverseproxy.service
[Unit]
Description=nmosreverseproxy
Wants=network-online.target
After=network.target network-online.target

[Service]
ExecStart=/usr/bin/proxylisting
User=ipstudio

[Install]
WantedBy=multi-user.target
Alias=ips-proxylisting.service
